# rootfs/etc/init.d/agetty-ttyS0
#!/sbin/openrc-run

name="agetty on ttyS0"
supervisor=supervise-daemon
command="/sbin/agetty"
command_args_foreground="--autologin root -s ttyS0 115200 xterm"
pidfile="/run/${RC_SVCNAME}.pid"
respawn_period=60

depend() {
	need mount-host9p
	after mount-host9p
}

start() {
	[ -c /dev/ttyS0 ] || return 0
	default_start
}

stop() {
	[ -c /dev/ttyS0 ] || return 0
	default_stop
}
