# rootfs/etc/init.d/agetty-hvc0
#!/sbin/openrc-run

name="agetty on hvc0"
supervisor=supervise-daemon
command="/sbin/agetty"
command_args_foreground="--autologin user42 -s hvc0 115200 xterm"
pidfile="/run/${RC_SVCNAME}.pid"
respawn_period=60

depend() {
	after agetty-ttyS0
}
