# rootfs/etc/init.d/mount-host9p
#!/sbin/openrc-run
# Mount virtio-9p export at /mnt/host (tag host9p). Retries until the device is up.

depend() {
	need localmount
}

start() {
	ebegin "Mount host9p at /mnt/host"

	modprobe 9pnet_virtio 2>/dev/null
	modprobe 9p 2>/dev/null

	_tries=0
	while [ "$_tries" -lt 60 ]; do
		if /usr/local/sbin/mount-host-share 2>/dev/null; then
			eend 0
			return 0
		fi
		sleep 1
		_tries=$((_tries + 1))
	done

	ewarn "host9p mount failed (no virtio-9p device yet)"
	eend 0
	return 0
}

stop() {
	need_stop=no
}
