# rootfs/etc/profile.d/10-locale-fr.sh
# Discover Linux — French locale for man(1) and system messages.
export LANG=fr_FR.UTF-8
export LC_MESSAGES=fr_FR.UTF-8
